# supervisor/windows.d/gettnship_login.ini
[program:gettnship_login]
directory=C:\Users\admin\PycharmProjects\sf_spider\gettnship  # 本地项目路径
command=C:\app\anaconda3\envs\sf_spider\python.exe -m scrapy crawl gettnship_login -s DEBUG=True  # 本地虚拟环境Python路径
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
environment=PYTHONPATH=C:\Users\admin\PycharmProjects\sf_spider;C:\Users\admin\PycharmProjects\sf_spider\gettnship
stdout_logfile=C:\Users\admin\PycharmProjects\sf_spider\logs\login.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5