# supervisor/windows.d/gettnship_shipments.ini
[program:gettnship_shipments]
directory=C:\Users\admin\PycharmProjects\sf_spider\gettnship  # 进入项目目录
command=C:\app\anaconda3\envs\sf_spider\python.exe -m scrapy crawl gettnship_shipments -s DEBUG=True
environment=PYTHONPATH=C:\Users\admin\PycharmProjects\sf_spider  # 保留根目录到 Python 路径
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=C:\Users\admin\PycharmProjects\sf_spider\logs\shipments.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5