# supervisor/windows.d/gettnship_shipments.ini
[program:gettnship_shipments]
directory=C:\Users\admin\PycharmProjects\sf_spider\gettnship
command=C:\app\anaconda3\envs\sf_spider\python.exe -m scrapy crawl gettnship_shipments -s DEBUG=False
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
environment=PYTHONPATH=C:\Users\admin\PycharmProjects\sf_spider;C:\Users\admin\PycharmProjects\sf_spider\gettnship
stdout_logfile=C:\Users\admin\PycharmProjects\sf_spider\logs\gettnship_shipments.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5

