# supervisor/conf.d/gettnship_login.ini
[program:gettnship_login]
directory=/opt/sf_spider  # 项目根路径
command=/opt/sf_spider/.venv/bin/python -m scrapy crawl gettnship_login -s DEBUG=False
autostart=true
autorestart=true
startretries=3
user=root
redirect_stderr=true
environment=PYTHONPATH=/opt/sf_spider
stdout_logfile=/opt/sf_spider/logs/login.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5