# supervisor/conf.d/scrapy_spider.ini
[program:gettnship_spider]
directory=/opt/sf_spider  # 项目路径
command=/opt/sf_spider/venv/bin/python -m scrapy crawl gettnship_shipments -s DEBUG=False  # 启动命令
autostart=true  # 系统启动时自动启动
autorestart=true  # 进程意外退出时自动重启
startretries=3  # 启动失败时的重试次数
user=root  # 运行用户
redirect_stderr=true  # 错误输出重定向到标准输出
environment=PYTHONPATH=/opt/sf_spider  # 环境变量设置
stdout_logfile=/opt/sf_spider/logs/spider.log  # 日志文件路径
stdout_logfile_maxbytes=10MB  # 单个日志文件最大大小
stdout_logfile_backups=5  # 保留的旧日志文件数量